# tweet.service
[Unit]
Description=Tweet current IP

[Service]
User=replaceusername
ExecStart=replaceexecpath

[Install]
WantedBy=multi-user.target
