# tweet.service
[Unit]
Description=Tweet current IP
After=ntpdate.service

[Service]
User=replaceusername
ExecStart=replaceexecpath

[Install]
WantedBy=multi-user.target
